# packaging/deb/prerm
#!/usr/bin/env sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    systemctl --global disable keypilot.service >/dev/null 2>&1 || true
fi

exit 0
